# etc/slurmd.service.in
[Unit]
Description=Slurm node daemon
After=network.target
ConditionPathExists=@sysconfdir@/slurm.conf

[Service]
Type=forking
ExecStart=@sbindir@/slurmd
PIDFile=/var/run/slurmd.pid

[Install]
WantedBy=multi-user.target
