# etc/slurmdbd.service.in
[Unit]
Description=Slurm DBD accounting daemon
After=network.target
ConditionPathExists=@sysconfdir@/slurmdbd.conf

[Service]
Type=forking
ExecStart=@sbindir@/slurmdbd
PIDFile=/var/run/slurmdbd.pid

[Install]
WantedBy=multi-user.target
