# etc/slurmdbd.service.in
[Unit]
Description=Slurm DBD accounting daemon
After=network.target munge.service
ConditionPathExists=@sysconfdir@/slurmdbd.conf

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/slurmdbd
ExecStart=@sbindir@/slurmdbd $SLURMDBD_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/var/run/slurmdbd.pid
LimitNOFILE=65536
@SYSTEMD_TASKSMAX_OPTION@

[Install]
WantedBy=multi-user.target
