# etc/slurmdbd.service.in
[Unit]
Description=Slurm DBD accounting daemon
After=network.target munge.service
ConditionPathExists=@sysconfdir@/slurmdbd.conf

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/slurmdbd
ExecStart=@sbindir@/slurmdbd -D $SLURMDBD_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
LimitNOFILE=65536
@SYSTEMD_TASKSMAX_OPTION@

[Install]
WantedBy=multi-user.target
